# anstele/roles/maldet/files/etc/cron.d/maldet_www
5 2 * * * root /usr/local/maldetect/maldet -b -r /var/www/hosting 1 >> /dev/null 2>&1

