# anstele/roles/php/templates/etc/cron.d/chrootphp
*/5 *     * * *     root   /usr/local/sbin/chroot-php

