# anstele/roles/cronurl/templates/etc/init/cronurl.conf
start on (local-filesystems and net-device-up)
stop on runlevel [!2345]

exec {{ virtualenv }}/bin/python {{ virtualenv }}/src/cronurl/cronurl/server.py -d
respawn

